# packaging/postinst
#!/bin/sh
set -e

# Reload systemd manager configuration
systemctl daemon-reload

# Enable seed service to start on boot
systemctl enable seed.service

# Start seed service
systemctl start seed.service

exit 0
