# packaging/postinst
#!/bin/sh
set -e

# Reload systemd manager configuration
systemctl daemon-reload

# Enable netscope service to start on boot
systemctl enable netscope.service

# Start netscope service
systemctl start netscope.service

exit 0
