# packaging/seed.service
[Unit]
Description=The Seed - Mustard Seed Networks
Documentation=https://github.com/krisarmstrong/seed
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=seed
Group=seed
WorkingDirectory=/usr/local/seed
ExecStart=/usr/local/bin/seed --config /usr/local/seed/configs/seed.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
# Grant network capabilities for ICMP ping, ARP scanning, and packet capture
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/usr/local/seed/configs /usr/local/seed/logs
PrivateTmp=true

# Environment
Environment=HOME=/usr/local/seed

[Install]
WantedBy=multi-user.target
