# packaging/seed.service
[Unit]
Description=The Seed - Network Diagnostic Tool
Documentation=https://github.com/krisarmstrong/seed
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=seed
Group=seed
WorkingDirectory=/var/lib/seed
ExecStart=/usr/bin/seed
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
# Grant network capabilities for ICMP ping, ARP scanning, packet capture, and link config
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=yes

# Filesystem protection
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/seed /var/lib/seed /var/log/seed
PrivateTmp=yes

# Additional hardening
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
PrivateDevices=no
# Note: PrivateDevices=no required for network interface access

# Environment
Environment=HOME=/var/lib/seed
Environment=SEED_DATA_DIR=/var/lib/seed

[Install]
WantedBy=multi-user.target
